# configs/dev/supervisor_gunicorn_bokaru.conf
[program:gunicorn]
command=sh /usr/local/bin/gunicorn_bokaru.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/gunicorn.log
redirect_stderr=true
