# configs/prod/supervisor_gunicorn_bokaru.conf
[program:gunicorn]
command=sh /usr/local/bin/gunicorn_bokaru.sh
autostart=true
autorestart=true
stdout_logfile=/usr/local/bokaru/logs/gunicorn.log
redirect_stderr=true
